# deploy/phase_e/com.modelatlas.phase_e.macpro.3.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.modelatlas.phase_e.macpro.3</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/python3</string>
        <string>/Users/squishfam/phase_e/phase_e_worker.py</string>
        <string>--input</string>
        <string>/Users/squishfam/phase_e/shard_3.jsonl</string>
        <string>--output</string>
        <string>/Users/squishfam/phase_e/results_3.jsonl</string>
        <string>--resume</string>
        <string>--searxng</string>
        <string>http://100.64.0.1:8888</string>
        <string>--snippets-only</string>
        <string>--model</string>
        <string>qwen2.5:3b</string>
        <string>--banks</string>
        <string>CAPABILITY,QUALITY</string>
        <string>--timeout</string>
        <string>300</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/squishfam/phase_e</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>Crashed</key>
        <true/>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    <key>Nice</key>
    <integer>20</integer>
    <key>LowPriorityIO</key>
    <true/>
    <key>ProcessType</key>
    <string>Background</string>
    <key>StandardOutPath</key>
    <string>/tmp/modelatlas-phase_e-macpro-3.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/modelatlas-phase_e-macpro-3.err</string>
</dict>
</plist>
